# Library/LaunchAgents/ooo.o0-o.env.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ooo.o0-o.env</string>
    <!-- Export env vars into the GUI launchd session so app-bundle
         launches (not just shells) see them. Shells get these from
         sh/profile.d; this covers Finder/Dock/Spotlight launches.
         Add more with '; launchctl setenv KEY "VALUE"'. -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <string>-c</string>
        <string>launchctl setenv CLAUDE_CONFIG_DIR "{{ .chezmoi.homeDir }}/.config/claude"</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
